# hidd/hid.service
[Unit]
Description=HID daemon :)

[Service]
ExecStart=/usr/local/bin/hid

[Install]
WantedBy=multi-user.target
